# sv/mpdscribble/run
#!/bin/sh
exec ~/.local/bin/mpdscribble --no-daemon
