# sv/tmux-rtorrent/run
#!/bin/sh
set -e

tmux new-session -s rt -n rtorrent -d rtorrent
tmux new-window -at rt -n rssrtorrent rssrtorrent.py
exec chpst -b tmux-rtorrent pause
