# sv/ipfs/run
#!/bin/sh
exec /home/kitty/prog/src/go/bin/ipfs daemon
